# xkeysnail.service
[Unit]
Description=xkeysnail

[Service]
ExecStart=/usr/bin/xkeysnail --quiet %h/.xkeysnail/config.py
Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target
